# hauzer_utility_exporter/rootfs/etc/services.d/hauzer-utility-exporter/run
#!/usr/bin/with-contenv bashio

export PYTHONPATH=/usr/lib
exec python3 -m hauzer_utility_exporter
